# deploy/supervisor-retc-cts-worker.conf
; Supervisor config for the RETC-CTS queue worker.
;
; `php artisan queue:work` must run as a long-lived background process —
; running it by hand in an SSH session dies the moment that session closes.
; Supervisor keeps it alive and restarts it automatically if it crashes.
;
; Install (on a VPS/dedicated server with root access — NOT typical on
; shared cPanel hosting, see the Deployment Checklist for the cPanel
; alternative using a cron-driven `queue:work --stop-when-empty`):
;
;   sudo apt install supervisor          # Debian/Ubuntu
;   sudo cp deploy/supervisor-retc-cts-worker.conf /etc/supervisor/conf.d/
;   # edit the two "REPLACE_ME" paths below first
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start retc-cts-worker:*

[program:retc-cts-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /REPLACE_ME/path/to/retc-cts/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/REPLACE_ME/path/to/retc-cts
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=REPLACE_ME_deploy_user
numprocs=1
redirect_stderr=true
stdout_logfile=/REPLACE_ME/path/to/retc-cts/storage/logs/worker.log
stopwaitsecs=3600
